# systemd/gssuserproxy.service.in
[Unit]
Description=GSS User Proxy
Documentation=man:gssproxy(8)

[Service]
Type=notify
StandardError=journal
ExecStart=@sbindir@/gssproxy -i -u
Restart=on-failure
